# src/bin/serv.service
[Unit]
Description=SERV daemon
After=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=1
User=serviceserv
WorkingDirectory=/home/serviceserv/serv
ExecStart=/home/serviceserv/serv/serv -c serv.conf
SyslogIdentifier=serv
StartLimitInterval=0
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target