# src/bin/servnode.service
[Unit]
Description=My SERV Node
After=network.target

[Service]
Type=simple
User=servuser
WorkingDirectory=/home/servuser/serv
ExecStart=/home/servuser/serv/cosmovisor start
Restart=on-failure
StartLimitInterval=0
RestartSec=3
LimitNOFILE=65535
LimitMEMLOCK=209715200

[Install]
WantedBy=multi-user.target